# docker/configs/proxy/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://proxy-oji01.jpn.hp.com:8080
